# systemd/kiwix-wikipedia.service
[Unit]
Description=Local Kiwix Wikipedia Server
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/cadmus/Projects/lm-studio-discord-bot
Environment=LD_LIBRARY_PATH=/home/cadmus/Projects/lm-studio-discord-bot/vendor/kiwix-tools/usr/lib/aarch64-linux-gnu
ExecStart=/home/cadmus/Projects/lm-studio-discord-bot/vendor/kiwix-tools/usr/bin/kiwix-serve -i 127.0.0.1 -p 8090 /home/cadmus/Projects/lm-studio-discord-bot/data/wikipedia/wikipedia_en-simple_all_nopic_2026-06.zim
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
