# systemd/lm-studio-discord-bot.service
[Unit]
Description=LM Studio Discord Bot
Requires=lm-studio-qwen35-08b.service
Wants=kiwix-wikipedia.service
After=network-online.target lm-studio-qwen35-08b.service kiwix-wikipedia.service
StartLimitIntervalSec=5min
StartLimitBurst=10
ConditionPathExists=/home/cadmus/Projects/lm-studio-discord-bot/.env

[Service]
Type=simple
WorkingDirectory=/home/cadmus/Projects/lm-studio-discord-bot
ExecStart=/home/cadmus/Projects/lm-studio-discord-bot/.venv/bin/python /home/cadmus/Projects/lm-studio-discord-bot/bot.py
Restart=always
RestartSec=5
TimeoutStopSec=30
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
