# systemd/lm-studio-qwen35-08b.service
[Unit]
Description=Ensure LM Studio qwen3.5-0.8b Model
Wants=network-online.target
Requires=lm-studio-app.service
After=network-online.target lm-studio-app.service
StartLimitIntervalSec=10min
StartLimitBurst=4

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/cadmus/Projects/lm-studio-discord-bot
ExecStart=/home/cadmus/Projects/lm-studio-discord-bot/scripts/ensure_lm_studio_model.sh
TimeoutStartSec=15min

[Install]
WantedBy=default.target
